# addon-immich-kiosk-gallery/rootfs/etc/s6-overlay/s6-rc.d/immich-kiosk/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Immich Kiosk Gallery
# Handle service finish
# ==============================================================================

bashio::log.info "Immich Kiosk Gallery service finished"
