# addon-immich-kiosk-gallery/rootfs/etc/s6-overlay/s6-rc.d/immich-kiosk/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Immich Kiosk Gallery
# Runs the Immich Kiosk Gallery web server
# ==============================================================================

bashio::log.info "Starting Immich Kiosk Gallery..."

# Dependencies are handled by s6-rc system automatically
# No need for manual s6-svwait since we have dependencies.d

# Set executable permissions
chmod +x /usr/bin/immich_kiosk.py

# Start the web server
exec python3 /usr/bin/immich_kiosk.py
