# addon-immich-kiosk-gallery/rootfs/etc/s6-overlay/s6-rc.d/immich-kiosk/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Immich Kiosk Gallery
# Runs the Immich Kiosk Gallery web server
# ==============================================================================

bashio::log.info "Starting Immich Kiosk Gallery..."


# Set executable permissions
chmod +x /usr/bin/immich_kiosk.py

# Start the web server
exec python3 /usr/bin/immich_kiosk.py
