# addon-immich-kiosk-gallery/rootfs/etc/s6-overlay/s6-rc.d/init-immich-kiosk/up
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Immich Kiosk Gallery
# Initialization script
# ==============================================================================

bashio::log.info "Initializing Immich Kiosk Gallery..."

# Create necessary directories
mkdir -p /usr/share/immich-kiosk/templates
mkdir -p /usr/share/immich-kiosk/static/css
mkdir -p /usr/share/immich-kiosk/static/js

# Log configuration
bashio::log.info "Immich URL: $(bashio::config 'immich_url')"
bashio::log.info "Show memories: $(bashio::config 'immich_show_memories')"
bashio::log.info "Show albums: $(bashio::config 'immich_show_albums')"

bashio::log.info "Immich Kiosk Gallery initialization completed"
